# recipes-apps/app/files/helloworld.service
[Unit]
Description=A small hello world application

[Service]
Type=forking
ExecStart=/usr/bin/helloworld

[Install]
WantedBy=multi-user.target
